# supervisor_linux.ini
[inet_http_server]
port = :9001

[supervisorctl]
serverurl=http://10.79.2.13:9001/



[program:admin-api]
directory = /data/apps/saturn/admin
command = ./admin-api
stdout_logfile = /data/apps/saturn/admin/log/admin-server.log
stderr_logfile = /data/apps/saturn/admin/log/admin-server.log
autorestart=true
autostart=false

[program:linx-test]
directory = /export/home/maguozhi/code/go/test_prometheus
command = ./test

[program:ddc]
directory = /export/home/maguozhi/unrealddc
command = ./unreal_httpddc_server
stdout_logfile = /export/home/maguozhi/unrealddc/logs/console.log
stderr_logfile = /export/home/maguozhi/unrealddc/logs/console.log

[program:env-test]
command = ./env-test.py
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr


[supervisord]
logfile=./var/supervisord.log ; 主 log 文件
logfile_maxbytes=50MB        ; log 文件最大尺寸
logfile_backups=10           ; 保留的旧 log 文件数量
loglevel=trace                ; 日志级别
pidfile=./var/supervisord.pid ; PID 文件路径
nodaemon=false               ; 作为后台守护进程运行
minfds=1024                  ; 可用文件描述符的最小数量
minprocs=200                 ; 可用进程的最小数量
